# etc/systemd/system/autosuspend.service
[Unit]
Description=Autosuspend script

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/autosuspend.sh
