# etc/systemd/system/autosuspend.service
[Unit]
Description=Run autosuspend.sh

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/autosuspend.sh
